# meta-example/recipes-webui/iot2050-web-gateway-nginx/files/cockpit.service.d/no-tls.conf
[Service]
# Run cockpit on plain HTTP only (no internal TLS).
# TLS is terminated exclusively by nginx (single entry point).
# Avoids double encryption and conflicting certificates.
# Must be paired with --for-tls-proxy so cockpit still treats traffic as HTTPS.
ExecStartPre=
ExecStart=
ExecStart=/usr/lib/cockpit/cockpit-tls --no-tls
