# meta-example/recipes-webui/iot2050-web-gateway-nginx/files/cockpit.socket.d/loopback.conf
[Socket]
# Bind cockpit to loopback only (127.0.0.1).
# Prevents direct external access and enforces nginx as the only entry point.
ListenStream=
ListenStream=127.0.0.1:9090
