# meta-example/recipes-webui/iot2050-web-gateway-nginx/files/cockpit-wsinstance-http.service.d/for-tls-proxy.conf
[Service]
# Enable proxy-aware mode for cockpit behind nginx TLS termination.
# Browser → HTTPS (nginx) → HTTP (cockpit)
# Without this, cockpit misinterprets requests as HTTP and rejects
# websocket/origin checks after login.
ExecStart=
ExecStart=/usr/lib/cockpit/cockpit-ws --for-tls-proxy --port=0
