# meta-example/recipes-webui/iot2050-web-gateway-nginx/files/nginx.service.d/iot2050-web-gateway.conf
[Service]
ExecStartPre=
# Prepare web-gateway environment before starting nginx.
# Typically ensures required files (certificates, configs, symlinks) exist.
ExecStartPre=/usr/lib/iot2050/web-gateway/iot2050-web-gateway-prepare.sh
