# meta-example/recipes-webui/iot2050-cockpit-customization/files/cockpit.socket.d/override.conf
[Socket]
ExecStartPost=
ExecStartPost=-/usr/lib/iot2050/cockpit/update-issue
