# meta-example/recipes-webui/iot2050-cockpit-customization/files/cockpit-issue.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/lib/iot2050/cockpit/update-issue
